# supervisor/kafka2.conf
[program:kafka2]
command=/usr/bin/start-kafka.sh server2.properties /etc/kafka2/jks
autostart=true
autorestart=true