# supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh server.properties /etc/kafka/jks
autostart=true
autorestart=true